# www/josegpt.com/pkg/josegpt_comd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/josegpt.com"
daemon_logger="daemon.info"
daemon_user="_josegpt.com"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
