# main/ptserver-uacme/crontab
0 0 * * * chronic /usr/bin/uacme-client
