# grpc/srth-account/pkg/srth_accountd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/srth-account"
daemon_user="_srth-account"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
